# examples/surge-tdx-prover/mkosi/default/mkosi.extra/usr/lib/systemd/system/taiko-client.service
[Unit]
Description=Taiko Client
After=runtime-init.service
Requires=runtime-init.service

[Service]
User=taiko-client
Group=eth
Restart=on-failure
ExecStart=/usr/bin/taiko-client

[Install]
WantedBy=default.target
