# examples/surge-tdx-prover/mkosi/default/mkosi.extra/usr/lib/systemd/system/tdxs.socket
[Unit]
Description=TDXS Socket
After=
Requires=

[Socket]
ListenStream=/var/tdxs.sock
SocketMode=0660
SocketUser=root
SocketGroup=tdx
Accept=false

[Install]
WantedBy=sockets.target
